# contrib/chissoku-forwarder.service
[Unit]
Description=Chissoku to Mackerel Forwarder
Documentation=https://github.com/nakataki17/co2-mackerel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=chissoku
Group=dialout

# 実行ファイルのパス（symlink 経由でバージョン管理）
# 設定は EnvironmentFile のみ（手順は docs/deployment.md）
ExecStart=/opt/chissoku-forwarder/current/forwarder

# 環境変数ファイル（MACKEREL_API_KEY, CHISSOKU_BIN など）
EnvironmentFile=/etc/chissoku-forwarder/forwarder.env

# 再起動設定
Restart=always
RestartSec=5

# セキュリティ設定
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/chissoku-forwarder
PrivateTmp=true

# ログ設定
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chissoku-forwarder

[Install]
WantedBy=multi-user.target
